# quartz/quartz_0.2_all/etc/cron.d/quartz
*/5 * * * * root python3 /usr/bin/quartz01.py
